# systemd/mimeclipd.service
[Unit]
Description=mimeclip clipboard history daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
ExecStart=%h/.cargo/bin/mimeclipd
Restart=on-failure
RestartSec=2
Environment=RUST_LOG=info

[Install]
WantedBy=graphical-session.target
